# examples/systemd/telegram-agent-bridge.service
[Unit]
Description=Telegram Agent Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Create this file from config.example.env and chmod 600 it.
EnvironmentFile=%h/.config/telegram-agent-bridge.env
# Services do not inherit your interactive shell PATH. Add the directory that
# contains codex, claude, aider, gemini, or your chosen agent command.
Environment=PATH=%h/.local/bin:%h/.npm-global/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin
WorkingDirectory=%h
ExecStart=/usr/bin/env python3 %h/telegram-agent-bridge/telegram_agent_bridge.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
