# examples/launchd/com.user.telegram-agent-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.telegram-agent-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOU/.local/bin/telegram-agent-bridge-run</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/YOU/.local/bin:/Users/YOU/.npm-global/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/telegram-agent-bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/telegram-agent-bridge.log</string>
</dict>
</plist>
